# examples/com.duckdns.updater.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  com.duckdns.updater.plist — schedule duckdns-update.py to run every 5 minutes.

  SETUP:
  1. Replace the two placeholder paths below:
       __PYTHON_PATH__   →  output of `which python3`
       __SCRIPT_PATH__   →  full path to duckdns-update.py
                            e.g. /Users/matt/jellyfin-yoto-bridge/scripts/duckdns-update.py

  2. Copy this file to ~/Library/LaunchAgents/com.duckdns.updater.plist

  3. Load it:
       launchctl load ~/Library/LaunchAgents/com.duckdns.updater.plist

  4. Verify it's running:
       launchctl list | grep duckdns

  TO DISABLE:
       launchctl unload ~/Library/LaunchAgents/com.duckdns.updater.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.duckdns.updater</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON_PATH__</string>
        <string>__SCRIPT_PATH__</string>
    </array>

    <!-- Run every 300 seconds (5 minutes) -->
    <key>StartInterval</key>
    <integer>300</integer>

    <!-- Run once immediately when loaded, then every 5 min -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Capture script output for debugging. Comment these out to silence. -->
    <key>StandardOutPath</key>
    <string>/tmp/duckdns-update.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/duckdns-update.err</string>
</dict>
</plist>
